# configuration/environments_scripts/central_reverse_proxy/files/etc/systemd/system/tmux-management-panel.service
[Unit]
Description=Starts up the tmux management panel
Requires=-.mount
After=-.mount  

[Install]
RequiredBy=multi-user.target

[Service]
Type=oenshot
RemainAfterExit=true
ExecStart=/usr/local/bin/tmuxManagementConsole.sh
