# configuration/environments_scripts/central_reverse_proxy/files/etc/systemd/system/wiki.service
[Unit]
Description=Starts up the wiki
Requires=-.mount
After=-.mount  

[Install]
RequiredBy=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/home/wiki/serve.sh
ExecStop=/home/wiki/serve.sh stop