# configuration/crontabs/crontab-switchoverArchive
* * * * *	export PATH=/bin:/usr/bin:/usr/local/bin; switchoverArchive.sh /etc/httpd/conf.d/000-macros.conf 2 9 