# configuration/crontabs/crontab-mongo-health-check
* * * * *  export PATH=/bin:/usr/bin:/sbin:/usr/sbin:/usr/local/bin;i notify-unhealthy-mongodb 2>/var/log/notify-unhealthy-mongodb.err >/var/log/notify-unhealthy-mongodb.out
