# configuration/crontabs/crontab-update-trac-trac-urls
10 12 * * *  export PATH=/bin:/usr/bin:/usr/local/bin; update-tractrac-urls-to-archive.sh "/root/code" >/var/log/update-tractrac-urls-to-archive.out 2>/var/log/update-tractrac-urls-to-archive.err